# deploy/ansible/roles/db_backup/templates/prosper-db-backup.timer.j2
[Unit]
Description=Daily Prosper DB backup

[Timer]
OnCalendar={{ backup_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
