# deploy/ansible/roles/prosper/templates/prosper.service.j2
[Unit]
Description=Prosper Expense Tracker
After=mariadb.service docker.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={{ service_user }}
# PROSPER_IMAGE selects the image tag to run; prosper-update.sh rewrites it on each rollout.
EnvironmentFile={{ image_env_file }}
ExecStartPre=-/usr/bin/docker rm -f prosper-fe
ExecStart=/usr/bin/docker run --rm --net host --name prosper-fe --env-file {{ env_file }} ${PROSPER_IMAGE}
StandardOutput=append:{{ log_dir }}/app.log
StandardError=append:{{ log_dir }}/app.log
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
