# deploy/ansible/roles/prosper/templates/prosper-update.timer.j2
[Unit]
Description=Prosper image update check

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
